# tools/late-report/systemd/late-report.timer
[Unit]
Description=Late Report Timer - Run daily at 11:00 MSK
Requires=late-report.service

[Timer]
OnCalendar=11:00 Europe/Moscow
# Каждый день в 11:00 по московскому времени
AccuracySec=5min
Persistent=true

[Install]
WantedBy=timers.target
